# deploy/oni-testnet.service
# oni-testnet.service - systemd service unit for oni Bitcoin testnet node
#
# Installation:
#   sudo cp oni-testnet.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable oni-testnet
#   sudo systemctl start oni-testnet

[Unit]
Description=oni Bitcoin Testnet Node
Documentation=https://github.com/your-org/oni
After=network-online.target
Wants=network-online.target
Requires=local-fs.target

[Service]
Type=simple
User=oni
Group=oni

WorkingDirectory=/opt/oni
ExecStart=/opt/oni/bin/oni --network testnet --data-dir /var/lib/oni-testnet
ExecStop=/bin/kill -SIGTERM $MAINPID

Restart=on-failure
RestartSec=10
RestartPreventExitStatus=0

LimitNOFILE=65536
LimitNPROC=4096

MemoryMax=4G
MemoryHigh=3G

TimeoutStartSec=300
TimeoutStopSec=120

Environment=ONI_NETWORK=testnet
Environment=ONI_DATA_DIR=/var/lib/oni-testnet
Environment=ONI_LOG_LEVEL=info
Environment=ONI_RPC_PORT=18332
Environment=ONI_P2P_PORT=18333
Environment=ONI_METRICS_PORT=19100

EnvironmentFile=-/etc/oni/oni-testnet.conf

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes

ReadWritePaths=/var/lib/oni-testnet
ReadWritePaths=/var/log/oni

RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

CapabilityBoundingSet=
AmbientCapabilities=

SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

[Install]
WantedBy=multi-user.target
